# ops-config/etc/systemd/system/snap.certbot.renew.timer
[Unit]
# Auto-generated, DO NOT EDIT
Description=Timer renew for snap application certbot.renew
Requires=snap-certbot-5371.mount
After=snap-certbot-5371.mount
X-Snappy=yes

[Timer]
Unit=snap.certbot.renew.service
OnCalendar=*-*-* 03:27
OnCalendar=*-*-* 21:17

[Install]
WantedBy=timers.target
