# config-backup/var/spool/cron/gort1975.cron
# reboot at 2am every day
0 2 * * * sudo /usr/sbin/reboot now
#
# run startup script on reboot
@reboot sleep 10 && sh /home/gort1975/scripts/startup.sh >> /home/gort1975/scripts/cronlog.log 2>&1
#
# run maintenence tasks 3am daily
0 3 * * * sh /home/gort1975/snorkelology/tools/run-nightly-maintenance.sh >> /home/gort1975/scripts/cronlog.log 2>&1
#
# These are now run on backup server
# 0 12 * * * /home/gort1975/emailcheck.sh 
# 0 * * * * /home/gort1975/site-status.sh 
